# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:myapp]
command=/app/main